# ram-usage@.service
[Unit]
Description=Ram usage service

[Service]
ExecStart=ram-usage.sh
StandardInput=socket
StandardOutput=socket